# deploy/systemd/dano-token-refresh.timer
[Unit]
Description=Run Dano runtime token refresh every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
Persistent=true
Unit=dano-token-refresh.service

[Install]
WantedBy=timers.target

